# extra/init/nagios-ehf.conf
# nagios-ehf 
#
#
description "Nagios EHF"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

console log 
expect fork
exec /usr/local/nagios-ehf/sbin/nagios-ehf -c /usr/local/nagios-ehf/etc/nagios-ehf.cfg -d -p
